# .config/systemd/user/mako.path
[Path]
PathModified=%h/.config/mako/config
Unit=mako-restarter.service

[Install]
WantedBy=graphical-session.target
